# deploy/systemd/gr-dis-bridge.service
[Unit]
Description=gr-dis DIS bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gr-dis
Group=gr-dis
WorkingDirectory=/opt/gr-dis
ExecStart=/opt/gr-dis/.venv/bin/gr-dis bridge --config /etc/gr-dis/config.yaml
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s
KillMode=mixed
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gr-dis-bridge
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
LockPersonality=yes
RestrictRealtime=yes
# AF_UNIX required by libzmq self-pipe
RestrictAddressFamilies=AF_INET AF_UNIX
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
